# opt/tor-socks-farm/systemd/rotate.timer
[Unit]
Description=Periodic rotation timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
Unit=rotate.service

[Install]
WantedBy=timers.target
