# opt/tor-socks-farm/systemd/rotate.service
[Unit]
Description=Rotate Tor identities

[Service]
Type=oneshot
ExecStart=/opt/tor-socks-farm/scripts/rotate.sh
